# resources/etc/profile.d/docker.sh
#!/usr/bin/bash
if [[ $(groups) == *"docker"* ]]
then
 source /etc/environment.d/docker
fi